# infra/systemd/matrix-calendar-inbound.service
[Unit]
Description=Infinite Matrix Calendar Inbound-Sync
After=network.target docker.service
Documentation=https://github.com/LEVCON-AT/Infinite-Matrix

[Service]
Type=simple
User=matrix-bridge
Group=matrix-bridge
WorkingDirectory=/opt/matrix-calendar-inbound
EnvironmentFile=/opt/matrix-calendar-inbound/.env
ExecStart=/usr/bin/node /opt/matrix-calendar-inbound/dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=matrix-calendar-inbound

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ReadWritePaths=/opt/matrix-calendar-inbound/data
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# MemoryDenyWriteExecute=true  # deaktiviert: V8-JIT braucht W+X (siehe Memory project_vps_deploy_lessons)
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0027

# Resources
LimitNOFILE=65536
MemoryMax=384M

[Install]
WantedBy=multi-user.target
